# deploy/launchd/cc.studiosaas.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>cc.studiosaas.backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_ROOT__/.venv/bin/python</string>
        <string>__PROJECT_ROOT__/backend/scripts/backup_postgres.py</string>
        <string>backup</string>
        <string>--keep</string>
        <string>14</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>STUDIOSAAS_DATABASE_URL</key>
        <string>postgresql://__USER__@localhost:5432/studiosaas_local_test</string>
        <key>PATH</key>
        <string>/opt/homebrew/opt/postgresql@18/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__HOME__/.studiosaas/backup.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.studiosaas/backup.log</string>
</dict>
</plist>
